# deploy/systemd/gc-autopull-weekly.timer
[Unit]
Description=Dugout GC autopull weekly report (Sunday 06:00 ET)

[Timer]
OnCalendar=Sun *-*-* 06:00:00 America/New_York
Persistent=true
Unit=gc-autopull-weekly.service

[Install]
WantedBy=timers.target
