# deploy/systemd/gc-autopull-weekly.service
[Unit]
Description=Dugout GC autopull weekly self-report
After=network-online.target

[Service]
Type=oneshot
User=joelycannoli
WorkingDirectory=/home/joelycannoli/dugout
EnvironmentFile=/home/joelycannoli/dugout/.env
ExecStart=/usr/bin/env python3 -m tools.autopull.weekly_report --days 7
StandardOutput=append:/home/joelycannoli/dugout/logs/autopull-weekly.log
StandardError=append:/home/joelycannoli/dugout/logs/autopull-weekly.log
TimeoutStartSec=120
